# pam_oktad/pam_oktad.service.in
[Unit]
Description=Okta PAM daemon
Documentation=man:@progname@(8)
After=network.target

[Service]
Type=exec
ExecStart=@prefix@/@sbindir@/@progname@ -S
KillMode=process
Restart=on-failure
RestartSec=41s

[Install]
WantedBy=multi-user.target
