# pam_oktad/pam_oktad.service.in
[Unit]
Description=Okta PAM daemon
Documentation=man:@progname@(8)
After=network.target

[Service]
Type=exec
ExecStartPre=/usr/bin/install -d -o root -g root -m 0700 --context=system_u:object_r:var_auth_t:s0 /var/run/pam_okta
ExecStart=@prefix@/@sbindir@/@progname@ -S
KillMode=process
Restart=on-failure
RestartSec=41s

[Install]
WantedBy=multi-user.target
